# srcs/grafana/supervisord.conf
[supervisord]
nodaemon=true
[program:grafana]
command= sh -c "./bin/grafana-server && kill -s SIGTERM $(cat supervisord.pid)"
